# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:hermes-webui]
command=python /opt/hermes-webui/server.py
directory=/opt/data/workspaces/default
user=hermes
environment=HERMES_HOME="/opt/data",HOME="/opt/data",HERMES_WEBUI_HOST="0.0.0.0",HERMES_WEBUI_PORT="8787",HERMES_WEBUI_AGENT_DIR="/opt/hermes",HERMES_WEBUI_STATE_DIR="/opt/data/webui",HERMES_WEBUI_DEFAULT_WORKSPACE="/opt/data/workspaces/default",TERMINAL_CWD="/opt/data/workspaces/default",VIRTUAL_ENV="/opt/hermes/.venv",PATH="/opt/hermes/.venv/bin:/usr/local/bin:/usr/bin:/bin"
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:hermes-dashboard]
command=/opt/hermes/.venv/bin/hermes dashboard --host 0.0.0.0
directory=/opt/data/workspaces/default
user=hermes
environment=HERMES_HOME="/opt/data",HOME="/opt/data",VIRTUAL_ENV="/opt/hermes/.venv",PATH="/opt/hermes/.venv/bin:/usr/local/bin:/usr/bin:/bin",HERMES_DASHBOARD_BASIC_AUTH_USERNAME="%(ENV_HERMES_DASHBOARD_BASIC_AUTH_USERNAME)s",HERMES_DASHBOARD_BASIC_AUTH_PASSWORD="%(ENV_HERMES_DASHBOARD_BASIC_AUTH_PASSWORD)s",HERMES_DASHBOARD_BASIC_AUTH_SECRET="%(ENV_HERMES_DASHBOARD_BASIC_AUTH_SECRET)s"
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/supervisor/hermes-dashboard.log
stderr_logfile=/var/log/supervisor/hermes-dashboard-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:hermes-gateway]
command=/opt/hermes/.venv/bin/hermes gateway run
directory=/opt/data/workspaces/default
user=hermes
environment=HERMES_HOME="/opt/data",HOME="/opt/data",API_SERVER_ENABLED="true",API_SERVER_HOST="0.0.0.0",GATEWAY_ALLOW_ALL_USERS="true",TERMINAL_CWD="/opt/data/workspaces/default",VIRTUAL_ENV="/opt/hermes/.venv",PATH="/opt/hermes/.venv/bin:/usr/local/bin:/usr/bin:/bin"
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/var/log/supervisor/hermes-gateway.log
stderr_logfile=/var/log/supervisor/hermes-gateway-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
